# celerydeploy/supervisord.conf
[unix_http_server]
file=%(here)s/supervisor.sock

[supervisord]
logfile=supervisord.log
pidfile=supervisord.pid

[program:celeryd]
command=celeryd --loglevel=DEBUG
stdout_logfile=celeryd.log
stderr_logfile=celeryd.log

autostart=true
autorestart=false

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
